# xCAT/postscripts/xcatpostinit1
#!/bin/bash
# IBM(c) 2007 EPL license http://www.eclipse.org/legal/epl-v10.html
# This script is deprecated, please use xcatpostinit1.install instead
# chkconfig: 345 84 59
# description: service node postboot script hack
# processname: xcatpostinit

### BEGIN INIT INFO
# Provides: xcatpostinit
# Default-Start: 3 4 5
# Default-stop: 0 1 2 6
# Required-Start: $network $syslog
# Should-Start: cron
# Required-Stop:
# Short-Description: xCATpost
# Description: xCAT post boot script
### END INIT INFO

# Source function library.
if [ -x /etc/rc.d/init.d/functions ]; then
  . /etc/rc.d/init.d/functions
fi

case $1 in
restart)
  $0 stop
  $0 start
  ;;
status)
  echo -n "xcatpostinit1 runs only at boot, runs additional post scripts"
  ;;
stop)
  echo -n "nothing to stop "
  ;;
start)
  # check for the REBOOT specified in xcatinfo to run post boot scripts on reboot
  if [ -f /opt/xcat/xcatinfo ]; then
    REBOOT=`grep 'REBOOT' /opt/xcat/xcatinfo |cut -d= -f2`
  fi
  # if the xcatdsklspost file exists and this is a reboot - run xcatdsklspost with a mode of 6
  if [ "$REBOOT" = "TRUE" -a -r /opt/xcat/xcatdsklspost ]; then
        /opt/xcat/xcatdsklspost 6
  else
    # run /opt/xcat/xcatinstallpost
    if [ -r /opt/xcat/xcatinstallpost ]; then
      /opt/xcat/xcatinstallpost
    fi
  fi
  ;;
esac
