# xCATsn/debian/postinst
#!/bin/sh
# postinst script for openmpi
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
      # so conserver will start
      mkdir -p /var/log/consoles

      # remove any management node file
      if [ -f /etc/xCATMN ]; then
        rm  /etc/xCATMN
      fi

        ver=$(cat /etc/issue |awk '{print $2}')
        num=${ver%.*}
        file="xcat.conf.apach24"
        if [ $num -gt 12 ];then
            mv /etc/apache2/conf.d/xcat.conf.apach24 /etc/apache2/conf-enabled/
        fi


      apachedaemon='apache2'

      update-rc.d $apachedaemon enable

      if [ -f "/proc/cmdline" ]; then   # prevent running it during install into chroot image
       	 XCATROOT=/opt/xcat /etc/init.d/xcatd start
 	 /etc/init.d/$apachedaemon reload 
      fi
      echo "xCATsn is now installed"

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
