# xCAT-genesis-scripts/etc/udev/rules.d/98-mlx.rules
SUBSYSTEM=="pci", ATTRS{subsystem_vendor}=="15b3", ATTRS{subsystem_device}=="4010", RUN+="/sbin/loadmlxeth"
