# xCAT-genesis-scripts/etc/udev/rules.d/99-imm.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="04b3", ATTRS{idProduct}=="4010", RUN+="/sbin/setupimmnic"
