# xCAT-genesis-scripts/etc/udev/rules.d/99-tsm.rules
SUBSYSTEM=="net", ATTRS{idVendor}=="046b", ATTRS{idProduct}=="ffb0", RUN+="/sbin/setuptsmnic"
